# deploy/application.service.erb
[Unit]
Description = The simplest ruby applicaiton running with rackup

[Service]
Environment = PATH=<%= path %>
WorkingDirectory = <%= application_directory %>
Type = simple
ExecStart = <%= bundle_path %> exec rackup -o 0.0.0.0 config.ru

[Install]
WantedBy = multi-user.target

